# run_chat.conf
[program:run_chat]
directory = /var/www/
command = npm start --product
autostart = true
autorestart = true
stderr-logfile = /var/log/run_chat.log